# docker/containers/php/crontab
0 * * * * cd /var/www/company.birzha-leads.com && /usr/local/bin/php console check-companies
